# cookbooks/yudachi/files/etc/systemd/system/yudachi.timer
[Unit]
Description=timer to notify rainfall to slack

[Timer]
OnUnitActiveSec=10m
Persistent=true

[Install]
WantedBy=timers.target
