# cookbooks/hatenikki/files/etc/systemd/system/hatenikki.timer
[Unit]
Description=Publish diary with Hatenikki gem every morning

[Timer]
OnCalendar=*-*-* 10:00
Persistent=true

[Install]
WantedBy=timers.target
